# deploy/homesec-cam.service.j2
[Unit]
Description=HomeSec Cam Flask
After=network-online.target
Wants=network-online.target

[Service]
User={{ service_user }}
WorkingDirectory={{ app_dir }}
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 {{ app_dir }}/app.py
Restart=always
RestartSec=3
# Ensure camera devices are accessible
SupplementaryGroups=video

[Install]
WantedBy=multi-user.target
